# services/bootc/cron.d/bootc-fetch-apply-updates
#!/usr/bin/sh
if [ -e /run/ostree-booted ]
then
   /usr/bin/bootc upgrade --apply --quiet
fi
